# files/limiter-shd.service
[Unit]
Description=Limiter SHD Service
Documentation=FreenetLite
After=syslog.target network-online.target

[Service]
User=root
NoNewPrivileges=true
ExecStart=/usr/local/bin/notif-limit-shd

[Install]
WantedBy=multi-user.target